# systemd/adelie-dbus-bridge.service
[Unit]
Description=Adelie D-Bus bridge (per-user)
Documentation=https://github.com/adelie-ai/desktop-assistant/issues/281
# The bridge re-exposes the daemon's session-bus surface (the
# org.desktopAssistant.* object paths) but speaks to the daemon as an ordinary
# authenticated UDS client: it mints a JWT from adelie-mint, opens a UDS
# connection to the daemon, and forwards method calls + signals.
#
# Soft Wants=/After= (not BindsTo=/Requires=): since #316 the bridge reconnects
# to the daemon on its own — and re-mints a fresh token each time — so a brief
# daemon restart no longer needs the bridge restarted; it reconnects in place.
# The dependencies only express ordering + best-effort co-start.
Wants=desktop-assistant-daemon.service
After=desktop-assistant-daemon.service
Wants=adelie-mint.service
After=adelie-mint.service
After=default.target
PartOf=default.target
# Q4 (issue #281): restarting the bridge drops its single UDS connection, which
# cancels any D-Bus-initiated turn still in flight. This is accepted by design —
# D-Bus callers (the KDE plasmoid) fire-and-poll, and a turn's lifetime is bound
# to the bridge's connection, not to the one-shot caller's.
#
# Failure-mode shift (issue #281): before the cutover, daemon up ⇒ D-Bus up.
# With the bridge, a missing/crashed bridge leaves D-Bus clients (KDE) dark even
# while the daemon is healthy. Observe liveness with:
#   systemctl --user is-active adelie-dbus-bridge.service
#   busctl --user list | grep org.desktopAssistant

[Service]
Type=simple
ExecStart=%h/.cargo/bin/adelie-dbus-bridge
Restart=on-failure
RestartSec=2
# The bridge IS the live D-Bus surface (#318 name flip): it claims
# `org.desktopAssistant`, and the daemon no longer serves the in-process surface
# (DESKTOP_ASSISTANT_DBUS_INPROCESS=false by default — set it true + stop the
# bridge to revert). For a side-by-side QA instance, run a second bridge with
# ADELIE_BRIDGE_NAME=org.desktopAssistant.Bridge or .Dev.
Environment=ADELIE_BRIDGE_NAME=org.desktopAssistant
Environment=RUST_LOG=info,desktop_assistant_dbus_bridge=debug

[Install]
WantedBy=default.target
