# systemd/adelie-mint.service
[Unit]
Description=Adelie JWT minter (per-user)
Documentation=https://github.com/adelie-ai/desktop-assistant/issues/101
# Mints short-lived HS256 JWTs over a local UDS ($XDG_RUNTIME_DIR/adelie/mint.sock)
# for socket clients — the D-Bus bridge, and any UDS/WS client configured with a
# minter_socket (#316). Callers are authenticated by SO_PEERCRED; set
# ADELIE_MINT_GROUP to restrict minting to a Unix group on multi-user hosts.
#
# The signing key defaults to the same file the daemon uses (under
# $XDG_DATA_HOME/desktop-assistant/secrets), so minted tokens validate against
# the daemon out of the box. That file is provisioned by the daemon, hence the
# ordering below; hardening this secret's storage is tracked in #365.
After=desktop-assistant-daemon.service
After=default.target
PartOf=default.target

[Service]
Type=simple
ExecStart=%h/.cargo/bin/adelie-mint
Restart=on-failure
RestartSec=2
Environment=RUST_LOG=info

[Install]
WantedBy=default.target
