# systemd/desktop-assistant-daemon.service
[Unit]
Description=Desktop Assistant Backend Daemon
After=default.target

[Service]
Type=simple
ExecStart=%h/.cargo/bin/desktop-assistant-daemon
Restart=on-failure
RestartSec=2
Environment=RUST_LOG=info

[Install]
WantedBy=default.target
