# systemd/desktop-assistant-daemon-dev.service
[Unit]
Description=Desktop Assistant Backend Daemon (Development)
After=default.target

[Service]
# Type=simple: the dev daemon serves UDS only — there is no in-process D-Bus
# surface since the cutover (#319). For a dev D-Bus surface, run a side-name
# bridge: `ADELIE_BRIDGE_NAME=org.desktopAssistant.Dev adelie-dbus-bridge`.
Type=simple
ExecStart=%h/.cargo/bin/desktop-assistant-daemon
Restart=on-failure
RestartSec=2
Environment=RUST_LOG=info

[Install]
WantedBy=default.target
